# src/signer-server/install/sshgate-signer-server.service
# systemd unit for sshgate-signer-server (v2 scaffold).
#
# Installed by install/deploy.sh at /etc/systemd/system/sshgate-signer-server.service.
# The deploy script substitutes ${BIN}, ${USER}, ${API_KEY_FILE}, ${DB_PATH},
# and ${LISTEN_ADDR} before writing this file — see deploy.sh for the
# values it expects.
#
# Design notes:
#   - Runs under a dedicated, unprivileged user (User=). Never root.
#   - Hardening directives (NoNewPrivileges, ProtectSystem, etc.) match the
#     spirit of the v1 sshgate-signer-telegram unit; tighten further once v2.1 lands and
#     the full filesystem layout settles.
#   - Restart=on-failure with a 5s back-off gives operators a moment to
#     spot a crash loop without thrash-restarting on a real bug.
#   - WatchdogSec is intentionally unset for v2.0; we wire it in v2.1 when
#     the server learns to ping sd_notify.

[Unit]
Description=sshgate-signer-server (SSHGate v2 hosted approval service)
Documentation=file://__INSTALL_DIR__/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
Group=__USER__
ExecStart=__BIN__ \
  --api-key-file __API_KEY_FILE__ \
  --db __DB_PATH__ \
  --addr __LISTEN_ADDR__
Restart=on-failure
RestartSec=5

# Hardening — defence in depth for a daemon that owns approval state.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
ReadWritePaths=__STATE_DIR__

[Install]
WantedBy=multi-user.target
